# backup/systemd/openclaw-local-backup.timer
[Unit]
Description=Run OpenClaw local fallback backup every 6 hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
