# backup/systemd/openclaw-restic-backup.service
[Unit]
Description=OpenClaw encrypted restic backup

[Service]
Type=oneshot
ExecStart=%h/.openclaw/workspace/backup/run-restic-backup.sh
