# backup/systemd/openclaw-restic-backup.timer
[Unit]
Description=Run OpenClaw restic backup every 6 hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
