# backup/systemd/openclaw-local-backup.service
[Unit]
Description=OpenClaw local encrypted fallback backup

[Service]
Type=oneshot
ExecStart=%h/.openclaw/workspace/backup/run-local-restic-backup.sh
